# deploy/ansible/roles/timelapse/templates/timelapse.conf.j2
[program:timelapse]
command=
    {{ webapps_dir }}/pir_pub.py
directory=/tmp
autostart=false
autorestart=true
stdout_logfile=/var/log/iot.log
stderr_logfile=/var/log/iot.log
redirect_stderr=true